# 2 SEGUNDO/1 Cuatrimestre/Sistemas Operativos - SO/Practicas/Modulo 1/APUNTES_MOD1/crontab/cron2
* * * * * /home/manjaro/Desktop/APUNTES_MOD1/crontab/renice.sh 1 2> /home/manjaro/Desktop/APUNTES_MOD1/crontab/error.txt
