# 2 SEGUNDO/1 Cuatrimestre/Sistemas Operativos - SO/Practicas/Modulo 1/APUNTES_MOD1/crontab/cron1
* * * * *  /home/manjaro/Desktop/crontab/script.sh 2> /home/manjaro/Desktop/crontab/error.txt
