# 2 SEGUNDO/1 Cuatrimestre/Sistemas Operativos - SO/Practicas/Modulo 1/APUNTES_MOD1/crontab/script.sh
#!/bin/bash

v=1
echo $v >> /home/manjaro/Desktop/crontab/a.txt
